# fanafesa_back/crontab

* * * * * php /c/xampp/htdocs/tpvs_back/artisan schedule:run 
